# deployment/supervisor/celerybeatd.conf
[program:celerybeat]
environment=OGN_CONFIG_MODULE='config/default.py'
command=/var/www/html/ogn-python/venv/bin/celery -A ogn_python.collect beat -l info
directory=/var/www/html/ogn-python

user=www-data
numprocs=1
stdout_logfile=/var/log/celery/beat.log
stderr_logfile=/var/log/celery/beat.log
autostart=true
autorestart=true
startsecs=10

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true