# deployment/supervisor/redmon.conf
[program:redmon]
command=/usr/local/bin/redmon --base-path /redmon/
user=root
directory=/var/www/html
autostart=true
autorestart=true